# dci_agent/systemd/dci-agent.service
[Unit]
Description=DCI Agent

[Service]
Type=oneshot
ExecStart=/usr/bin/dci-agent
User=dci-agent

[Install]
WantedBy=default.target