# dci_agent/systemd/dci-agent.timer
[Unit]
Description=Trigger the DCI-Agent

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target